# linux/systemd/litellm.service
[Unit]
Description=LiteLLM — Anthropic-compatible gateway to Azure OpenAI / Vertex Gemini / etc.
After=network-online.target
Wants=network-online.target
# Tolerate slow first boots: Prisma migrations + engine-binary fetch can take
# 15–25s on a fresh DB, and 5 retries * 5s = 25s used to hit start-limit-hit.
StartLimitBurst=10
StartLimitIntervalSec=300

[Service]
Type=simple
WorkingDirectory=__APP_DIR__
EnvironmentFile=__ENV_FILE__
Environment=PATH=__PATH__
Environment=LITELLM_MODE=PRODUCTION
Environment=LITELLM_LOG=ERROR
Environment=LITELLM_LOCAL_ANTHROPIC_BETA_HEADERS=true
Environment=PRISMA_BINARY_CACHE_DIR=%h/.cache/prisma-python/binaries
Environment=PRISMA_HIDE_UPDATE_MESSAGE=1
Environment=PRISMA_SKIP_POSTINSTALL_GENERATE=1
ExecStart=__LITELLM_BIN__ --config __APP_DIR__/config.yaml --port __PORT__ --host 127.0.0.1 --num_workers 1 --telemetry False
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
